# integration_test/container.target
[Unit]
Description=Simplistic target for containers that avoids starting unnecessary systemd components